# packaging/scripts/postinstall
#!/bin/bash

# Set permissions for the helper tool
chown root:wheel /Library/PrivilegedHelperTools/slopn-helper
chmod 755 /Library/PrivilegedHelperTools/slopn-helper

# Set permissions for the launchd plist
chown root:wheel /Library/LaunchDaemons/com.webdunesurfer.slopn.helper.plist
chmod 644 /Library/LaunchDaemons/com.webdunesurfer.slopn.helper.plist

# Ensure the daemon is unloaded before loading again to force a restart
launchctl unload /Library/LaunchDaemons/com.webdunesurfer.slopn.helper.plist 2>/dev/null

# Ask user for Server Address and Token via AppleScript
# Note: This runs as root, so we use 'as user' logic if needed, but standard dialog usually works.
USER_SERVER=$(osascript -e 'display dialog "Enter SloPN Server Address (e.g. 1.2.3.4:4242):" default answer "" with title "SloPN Configuration" with icon note' -e 'text returned of result' 2>/dev/null)
USER_TOKEN=$(osascript -e 'display dialog "Enter SloPN Auth Token:" default answer "" with title "SloPN Configuration" with icon note with hidden answer' -e 'text returned of result' 2>/dev/null)
USER_SNI=$(osascript -e 'display dialog "Enter Mimic Target (SNI) [optional]:" default answer "www.google.com" with title "SloPN Configuration" with icon note' -e 'text returned of result' 2>/dev/null)

# Store these in a system-wide config if helper supports it, or just pass to GUI via a shared file
# For now, we will create a global config file that both Helper and GUI can read.
mkdir -p /Library/Application\ Support/SloPN
echo "{\"server\":\"$USER_SERVER\", \"token\":\"$USER_TOKEN\", \"sni\":\"$USER_SNI\", \"obfuscate\":true}" > "/Library/Application Support/SloPN/config.json"
chmod 644 "/Library/Application Support/SloPN/config.json"

# Create marker for GUI to know it's a fresh install/reconfig
touch "/Library/Application Support/SloPN/.new_install"
chmod 644 "/Library/Application Support/SloPN/.new_install"

# Load the daemon
launchctl load -w /Library/LaunchDaemons/com.webdunesurfer.slopn.helper.plist

# Generate a random IPC secret if it doesn't exist
SECRET_FILE="/Library/Application Support/SloPN/ipc.secret"
if [ ! -f "$SECRET_FILE" ]; then
    echo -n "$(openssl rand -hex 32)" > "$SECRET_FILE"
    chmod 644 "$SECRET_FILE" # Readable by root and users so GUI can see it
fi

exit 0
